# iio-sensor-proxy-libssc-openrc/iio-sensor-proxy-libssc.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/usr/lib/iio-sensor-proxy"

depend() {
	need dbus
}
